# root/etc/service/bitlbee/run
#!/bin/bash
chown -R app:app /var/lib/bitlbee
exec /usr/sbin/bitlbee -F -n
